# scripts/com.squad0.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.squad0.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/squad0/squad0</string>
        <string>start</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/opt/squad0</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/opt/squad0/data/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/opt/squad0/data/logs/launchd-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/squad0:/opt/homebrew/bin:/Users/james/.local/bin:/Users/james/go/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/james</string>
        <key>OLLAMA_HOST</key>
        <string>http://localhost:11434</string>
    </dict>
</dict>
</plist>
